# web/secure_file_store/supervisord-chromium.conf
[program:chrome]
command=/usr/bin/node /admin/client.js
user=admin
